# launchagents/com.mk.gog.gmail-personal-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.mk.gog.gmail-personal-serve</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>set -euo pipefail; source /Users/mk/.openclaw/secrets/secrets.env; exec /opt/homebrew/bin/gog --client personal --account kravch@gmail.com gmail watch serve --bind 127.0.0.1 --port 8788 --path / --token "$GMAIL_PERSONAL_PUSH_TOKEN" --hook-url http://127.0.0.1:18789/hooks/gmail-personal --hook-token "$OPENCLAW_HOOK_TOKEN" --include-body --max-bytes 20000</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <key>StandardOutPath</key><string>/Users/mk/Library/Logs/gog-gmail-personal-serve.log</string>
    <key>StandardErrorPath</key><string>/Users/mk/Library/Logs/gog-gmail-personal-serve.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
  </dict>
</plist>
